# infra/systemd/yf-reclassify-topics.timer
[Unit]
Description=YOUFEED — reclassify articoli 2 volte al giorno

[Timer]
# 04:00 e 16:00 UTC — fasce a basso traffico, due passate al giorno per
# riassorbire rapidamente le correzioni manuali di topic fatte via admin.
OnCalendar=*-*-* 04:00:00 UTC
OnCalendar=*-*-* 16:00:00 UTC
Persistent=true
Unit=yf-reclassify-topics.service

[Install]
WantedBy=timers.target
