# infra/systemd/yf-worker@.service
[Unit]
Description=YOUFEED RQ worker (%i)
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=youfeed
Group=youfeed
WorkingDirectory=/opt/youfeed/backend
EnvironmentFile=/opt/youfeed/.env
# %i è il nome della coda passato come argomento all'unit (es. yf-worker@fetch_rss.service)
ExecStart=/opt/youfeed/.venv/bin/rq worker --url ${RQ_REDIS_URL} %i

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/youfeed /var/log/youfeed
ProtectHome=true

[Install]
WantedBy=multi-user.target
